# runway.conf
[program:$RUNWAY_INSTALL_NAME]
command=$RUNWAY_INSTALL_BINDIR/$RUNWAY_INSTALL_NAME
stdout_logfile=/var/log/supervisor/$RUNWAY_INSTALL_NAME.log
redirect_stderr=true