# runway.conf
[program:runway]
command=/opt/bin/runway